# debian/static/etc/network/if-pre-up.d/quicktun
#!/bin/sh
test -n "${IF_QT_REMOTE_ADDRESS}" || exit 0
/usr/sbin/openvpn --mktun --dev "${IFACE}" --user quicktun
