# debian/static/etc/network/if-up.d/quicktun
#!/bin/sh
test -n "${IF_QT_REMOTE_ADDRESS}" || exit 0
/usr/bin/daemon -n "quicktun.${IFACE}" -u quicktun -i -l daemon.err -b daemon.debug -o daemon.debug /usr/sbin/quicktun.debian
